# systemd/munin-cron.service
[Unit]
Description=Survey monitored computers (munin)
After=network.target

[Service]
User=munin
ExecStart=/usr/bin/munin-cron
